# deploy/wsl/systemd/govpress-compose.service
[Unit]
Description=GovPress Docker Compose Stack
Wants=network-online.target docker.service
After=network-online.target docker.service docker.socket

[Service]
Type=oneshot
WorkingDirectory=/home/wavel/projects/GovPress_PDF_MD/deploy/wsl
Environment=HOME=/home/wavel
RemainAfterExit=yes
ExecStartPre=/bin/bash -lc "for i in $(seq 1 20); do docker info >/dev/null 2>&1 && exit 0; sleep 3; done; exit 1"
ExecStart=/bin/bash --noprofile --norc -c "/home/wavel/projects/GovPress_PDF_MD/deploy/wsl/bin/compose.sh up -d"
ExecStop=/bin/bash --noprofile --norc -c "/home/wavel/projects/GovPress_PDF_MD/deploy/wsl/bin/compose.sh down"
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
