# deploy/wsl/systemd/govpress-cloudflared.service
[Unit]
Description=GovPress Cloudflared Tunnel
Wants=network-online.target govpress-caddy.service
After=network-online.target govpress-caddy.service

[Service]
Type=simple
Environment=HOME=/home/wavel
Environment=PATH=/home/wavel/.local/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/home/wavel/projects/GovPress_PDF_MD/deploy/wsl
EnvironmentFile=/home/wavel/projects/GovPress_PDF_MD/deploy/wsl/.env
ExecStartPre=/bin/bash -lc 'test -n "$CLOUDFLARE_TUNNEL_TOKEN"'
ExecStart=/bin/bash -lc 'exec cloudflared tunnel --no-autoupdate --protocol http2 --url http://127.0.0.1:8080 run --token "$CLOUDFLARE_TUNNEL_TOKEN"'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
