# deploy/vps/systemd/govpress-compose.service
[Unit]
Description=GovPress Docker Compose Stack
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=__COMPOSE_DIR__
RemainAfterExit=yes
ExecStart=/usr/bin/docker compose -f __COMPOSE_DIR__/docker-compose.yml up -d
ExecStop=/usr/bin/docker compose -f __COMPOSE_DIR__/docker-compose.yml down
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
